# bird-webserver.service
[Unit]
Description=Smart Bird Feeder — Photo Browser & Health Dashboard
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/smart-bird-feeder
ExecStart=/home/pi/smart-bird-feeder/venv/bin/python3 webserver.py
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

MemoryMax=100M

[Install]
WantedBy=multi-user.target
